# salt/hwaas-service/supervisor.conf
[program:hwaas-site]
command=node .
directory=/home/hwaas/hwaas-site
user=hwaas
